# cron.daily/HF.nightly-build.cron
#!/bin/bash

. /etc/profile

/var/lib/HappyFace-Integration-Server/nightly-build.sh -r -e
